# services/linuxia-health-report.timer
[Unit]
Description=LinuxIA health report timer

[Timer]
OnCalendar=*-*-* 03:10:00
Persistent=true
Unit=linuxia-health-report.service

[Install]
WantedBy=timers.target
