# services/systemd/linuxia-healthcheck.timer
[Unit]
Description=Run LinuxIA healthcheck daily

[Timer]
OnCalendar=*-*-* 03:05:00
Persistent=true
Unit=linuxia-healthcheck.service

[Install]
WantedBy=timers.target
