# services/systemd/linuxia-repair.service
[Unit]
Description=LinuxIA auto-repair (VM100)
ConditionPathExists=/opt/linuxia/scripts/linuxia-repair.sh
After=network-online.target local-fs.target
Wants=network-online.target

[Service]
RemainAfterExit=yes
KillMode=process
Type=oneshot
ExecStart=/opt/linuxia/scripts/linuxia-repair.sh
TimeoutStartSec=300
